# autogit.target
[Unit]
Description=My autogit download target
Wants=autogit@testrepo.service

[Install]
WantedBy=default.target
