# autogit.target
[Unit]
Description=My autogit download target

[Install]
WantedBy=default.target
