# autogit.target
[Unit]
Description=My autogit download target
After=network.target

[Install]
WantedBy=default.target
